# guest-warning/deb-maintainer-scripts/postinst
#!/bin/bash

set -e

systemctl --global enable guest-users-guest-warning
