# guest-warning/deb-maintainer-scripts/prerm
#!/bin/bash

set -e

systemctl --global disable guest-users-guest-warning
